# scripts/starts/shelleasytier.service
[Unit]
Description=ShellEasytier Service
After=network.target syslog.target
Wants=network.target

[Service]
Type=simple
ExecStart=/etc/ShellEasytier/start.sh daemon-run
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
