# scripts/starts/shelleasytier-web.procd
#!/bin/sh /etc/rc.common

START=96
USE_PROCD=1

APPDIR=/etc/ShellEasytier

start_service() {
    procd_open_instance
    procd_set_param command "$APPDIR/start.sh" web-daemon-run
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    kill -TERM "$(cat /tmp/ShellEasytier/easytier-web.pid 2>/dev/null)" 2>/dev/null
    rm -f /tmp/ShellEasytier/easytier-web.pid
    killall easytier-web-embed 2>/dev/null
}
