# scripts/starts/shelleasytier.procd
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

APPDIR=/etc/ShellEasytier

start_service() {
    procd_open_instance
    procd_set_param command "$APPDIR/start.sh" daemon-run
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    "$APPDIR/start.sh" stop >/dev/null 2>&1
}
